# web/webapp.conf
[program:webapp]
command=/workspace/src/app/app
directory=/workspace/src/app
